# infra/systemd/vd-api.service
[Unit]
Description=Video Digestor - FastAPI Backend
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=vd
Group=vd
WorkingDirectory=/opt/vd/repo
EnvironmentFile=/opt/vd/.env
ExecStart=/usr/bin/bash -lc 'exec "${UV_PROJECT_ENVIRONMENT:-/var/cache/vd/project-venv}/bin/uvicorn" apps.api.app.main:app --host 127.0.0.1 --port 8000 --workers 2 --log-level info'
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vd-api

# Hardening
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
